# templates/hubot.upstart.erb
start on runlevel [2345]
stop on runlevel [!2345]

setuid hubot
setgid hubot

chdir <%= scope.lookupvar('hubot::root_dir') %>/<%= scope.lookupvar('hubot::bot_name') %>

env HOME=<%= scope.lookupvar('hubot::root_dir') %>
<% @exports.sort_by{|k,v| k}.each do |k,v| -%>
env <%= "#{k}=\"#{v}\"" %>
<% end -%>

pre-start script
  npm install
end script

exec <%= scope.lookupvar('hubot::root_dir') %>/<%= scope.lookupvar('hubot::bot_name') %>/node_modules/.bin/hubot --adapter <%= scope.lookupvar('hubot::adapter') %>
